# launchd/com.lihuan.morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lihuan.morning-brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/lihuan/.claude/scripts/run_morning_brief_codex.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>10</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>40</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/lihuan/Documents/学习提升/AI/claude/桌面 claude code/Scheduled/morning-brief-email/logs/cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lihuan/Documents/学习提升/AI/claude/桌面 claude code/Scheduled/morning-brief-email/logs/cron.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/lihuan/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/lihuan</string>
        <key>CODEX_HOME</key>
        <string>/Users/lihuan/.codex</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/lihuan</string>

    <key>ThrottleInterval</key>
    <integer>1</integer>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
