# launchd/com.lihuan.youtube-voc-monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lihuan.youtube-voc-monthly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/lihuan/Documents/跨境业务/美国亚马逊业务运营/social_media_voc/youtube_voc/run_monthly.sh</string>
    </array>

    <!-- Every Wednesday at 11:00 AM ET; script self-checks day ≤ 7 for first-Wednesday only -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>3</integer>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/youtube_voc_monthly_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/youtube_voc_monthly_launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
